# com.thanprojects.planterplanner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.thanprojects.planterplanner</string>

  <key>WatchPaths</key>
  <array>
    <string>/Users/thandrage/Claude/Projects/Planter Planner/.push</string>
  </array>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/thandrage/Claude/Projects/Planter Planner/push.sh</string>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/thandrage/Claude/Projects/Planter Planner/.push.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/thandrage/Claude/Projects/Planter Planner/.push.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
